# nekowm.service
[Unit]
Description="NekoWM startup service"

[Service]
Type=oneshot
ExecStart=/usr/bin/nekowm-login

[Install]
WantedBy=multi-user.target
